# auto_gitclone.cron
*/5 * * * * /usr/src/qbuy_test/autogitclone.sh


